# supervisor.conf
[program:timelapse]
command = python raspiLapseCam.py
directory = /home/pi
autostart = true
user = pi
stdout_logfile = /var/log/supervisor/timeline.log